# airootfs/etc/systemd/system/arCNiso-prepare.service
[Unit]
Description=arCNiso prepare

[Service]
Type=oneshot
ExecStart=/usr/local/arCNiso/prepare.sh

[Install]
WantedBy=multi-user.target
